# app/conf.d/code-server.conf
[program:code-server]
command=code-server --bind-addr 0.0.0.0:8008 --cert --auth none 
user=ubuntu
autorestart=true
stdout_logfile=/app/logs/code-server.log
redirect_stderr=true
#autostart=true
#stopwaitsecs=30
#stdout_logfile=/dev/stdout
#stdout_logfile_maxbytes=0
#stderr_logfile=/dev/stderr
#stderr_logfile_maxbytes=0
