# app/conf.d/xsession.conf
[program:xsession]
#command=/etc/X11/Xsession
command=jwm -display :9 -f /app/system.jwmrc
user=ubuntu
autorestart=true
stdout_logfile=/app/logs/xsession.log
redirect_stderr=true